# conf/systemd/mumble.service
[Unit]
Description = Mumble Client  
Requires = systemd-user-sessions.service network.target sound.target
After = multi-user.target

[Service]
User = talkkonnect  
Group = talkkonnect
Type = simple  
ExecStart = /home/talkkonnect/bin/talkiepi -server 192.168.1.30:64738 -password fablab -username zero01 -channel 7 -insecure
Restart = always 
RestartSec = 5

[Install]
WantedBy = multi-user.target
